# supervisor.conf
[supervisord]
nodaemon=true

[program:amuled]
command=/usr/bin/amuled

[program:amuleweb]
command=/usr/bin/amuleweb
